# custom-root/etc/skel/.config/autostart/guestpc.desktop
[Desktop Entry]
Name=GuestPC
Comment=GuestPC by Zimmerbuch.de
Exec=/usr/local/guestpc/guestpc-login.sh
Icon=preferences-desktop-remote-desktop
Terminal=false
Type=Application
X-GNOME-Autostart-Phase=Applications
